# config/supervisor/mailsink_supervisor.conf
; Mailsink configuration
[program:mailsink]
command=/srv/www/virtualenvs/mail_server/bin/mailsinkd -w 8025 -s 25
stopsignal=TERM
stderr_logfile_backups=5
