# config/supervisor/rookeries_supervisor.conf
; UWSGI for Rookeries
[program:rookeries_uwsgi]
user=www-data
command=/srv/www/virtualenvs/uwsgi/bin/uwsgi --ini /srv/www/config/uwsgi/rookeries-uwsgi.ini
stopsignal=TERM
stderr_logfile_backups = 10
